# supervisord.conf
[supervisord]
nodaemon=true

[program:spdo]
command=/usr/bin/pidproxy /var/interlegis/spdo/zinstance/var/spdo.pid /var/interlegis/spdo/zinstance/bin/instance start
